# systemd/marketplace-monitor.service
[Unit]
Description=Marketplace Monitor — dark-web breach feed dashboard
After=network-online.target tor.service
Wants=network-online.target

[Service]
Type=simple
User=marketplace-monitor
Group=marketplace-monitor
WorkingDirectory=/opt/marketplace-monitor
EnvironmentFile=/etc/marketplace-monitor.env
ExecStart=/opt/marketplace-monitor/.venv/bin/marketplace-monitor
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=marketplace-monitor

# Hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/opt/marketplace-monitor/data /opt/marketplace-monitor/config

[Install]
WantedBy=multi-user.target
